# systemd/mp3tag-backend.service
[Unit]
Description=MP3 Tag Editor - Backend
After=network-online.target remote-fs.target
Wants=network-online.target

[Service]
Type=simple
User=kenneth
WorkingDirectory=/home/kenneth/mp3-tag-editor/backend
ExecStart=/home/kenneth/mp3-tag-editor/backend/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

EnvironmentFile=/home/kenneth/mp3-tag-editor/backend/.env

[Install]
WantedBy=multi-user.target
